# root/etc/services.d/plexAutoscan/run
#!/usr/bin/python
echo "Starting Plex Media Server Auto Scanner."
exec /app/scan.py server --loglevel=INFO